# infra/systemd/veda-restore-drill.timer
[Unit]
Description=Run Veda PostgreSQL restore drill weekly

[Timer]
OnCalendar=Sun *-*-* 02:20:00
AccuracySec=5min
Persistent=true
Unit=veda-restore-drill.service

[Install]
WantedBy=timers.target
