# infra/systemd/veda-healthwatch.timer
[Unit]
Description=Run Veda Trading AI healthwatch every 2 minutes

[Timer]
OnBootSec=90s
OnUnitActiveSec=2min
AccuracySec=20s
Persistent=true
Unit=veda-healthwatch.service

[Install]
WantedBy=timers.target
